# debian/prerm
#!/bin/sh
# prerm script for n4d-ldap
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <prerm> `remove'
#        * <old-prerm> `upgrade' <new-version>
#        * <new-prerm> `failed-upgrade' <old-version>
#        * <conflictor's-prerm> `remove' `in-favour' <package> <new-version>
#        * <deconfigured's-prerm> `deconfigure' `in-favour'
#          <package-being-installed> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

PACKAGE_NAME="n4d-ldap"
ORIGINAL_FILE="/etc/init.d/slapd"
DIVERT_FILE="/usr/share/n4d-ldap/diversions/etc.initd.slapd.divert"
LLIUREX_FILE="/usr/sbin/slapd-daemon-starter"


case "$1" in
    upgrade|deconfigure)
    ;;

    remove|purge)

	if [ -f ${DIVERT_FILE} ] ; then
        	rm -rf $LLIUREX_FILE || true
		systemctl disable slapd || true
		rm -rf /lib/systemd/system/slapd.service || true
		systemctl daemon-reload || true
                dpkg-divert --package  ${PACKAGE_NAME} --rename --quiet --remove ${ORIGINAL_FILE}|| true
                systemctl daemon-reload || true
		update-rc.d slapd enable 2 3 4 5 || true
		
        fi


    ;;

    failed-upgrade)
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
